# deb/ma-agent.postinst
#!/bin/bash

set -e

if [ -x "/etc/init.d/ma-agent" ]; then
  update-rc.d ma-agent defaults >/dev/null 2>&1
fi
